# deploy/systemd/codex-telegram.service
[Unit]
Description=Codex Telegram Codex bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/.local/share/codex-telegram
Environment=PYTHONPATH=%h/.local/share/codex-telegram/src
ExecStart=/usr/bin/python3 -m codex_telegram.cli --env-file %h/.local/share/codex-telegram/.env
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
