# deploy/systemd/hermes-telegram.service
[Unit]
Description=Hermes Telegram Codex bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/Desktop/hermes-telegram
Environment=PYTHONPATH=%h/Desktop/hermes-telegram/src
ExecStart=/usr/bin/python3 -m hermes_telegram.cli --env-file %h/Desktop/hermes-telegram/.env
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
